# recipes-swtpm/swtpm-service/swtpm-service/swtpm.service
[Unit]
Description=SWTPM
Before=parsec.service

[Service]
User=parsec
Restart=always
RestartSec=5s
RemainAfterExit=yes
WorkingDirectory=/userdata
ExecStartPre=+mkdir -m 0700 -p /userdata/parsec
ExecStartPre=+chown parsec /userdata/parsec
ExecStart=/usr/bin/swtpm.sh

[Install]
WantedBy=default.target
