# recipes-swtpm/swtpm-service/swtpm-service/swtpm.service
[Unit]
Description=SWTPM
Before=parsec.service

[Service]
WorkingDirectory=/userdata
ExecStart=/bin/tpm_server

[Install]
WantedBy=default.target
